# lib/systemd/system/upgrade-archlinux.timer
[Unit]
Description=Daily upgrade of archlinux

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=multi-user.target
