# etc/upgrade-archlinux.service
[Unit]
Description=upgrade archlinux

[Service]
Type=oneshot
ExecStart=/bin/bash -c "pacman -Syu --noconfirm && pacman -Scc --noconfirm"
Nice=19
