# etc/yiend
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

SS_SERVER_ADDR=
SS_PASSWORD=

#depend() {
#
#}

start() {
	/opt/shadowsocks/shadowsocks/local.py -s $SS_SERVER_ADDR -k $SS_PASSWORD -d start
	/opt/alidns/src/main.py &>/tmp/alidns.log &
}

stop() {
	/opt/shadowsocks/shadowsocks/local.py -s $SS_SERVER_ADDR -d stop
}
